# batler.service
[Unit]
Description=Batler Server Remote Jobs Manager

[Service]
ExecStart=/usr/bin/batler-server
KillMode=mixed
Restart=always
IgnoreSIGPIPE=no